# crontab
PYTHONPATH=/home/greatagainkr

0 */6 * * * root /usr/local/bin/python /home/greatagainkr/greatagain_parser_naver/run.py > /proc/1/fd/1 2>/proc/1/fd/2
